# docker/supervisord.conf
; =============================================================================
; Supervisord — AI Dev Workspace Process Manager
; =============================================================================
; Manages 2 UI apps — OmniRoute and hermes-agent are started BEFORE supervisord
; by entrypoint.sh (OmniRoute first so hermes-agent has an LLM backend).
;
;   - Hermes Workspace  (port 3000)  — connects to hermes-agent on localhost:8642
;   - Aion UI           (port 3005)
;
; Started by entrypoint.sh BEFORE supervisord (not managed here):
;   - OmniRoute         (port 4000)  — LLM router (hermes-agent's LLM backend)
;   - hermes gateway    (port 8642)  — OpenAI-compatible + Hermes APIs
;   - hermes dashboard  (port 9119)  — Sessions, memory, skills, jobs APIs
; =============================================================================

[supervisord]
nodaemon=true
user=hermes
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/var/log/supervisor



; ---------------------------------------------------------------------------
; Hermes Workspace
; Connects to hermes-agent running on localhost inside the same container.
; ---------------------------------------------------------------------------
[program:hermes]
command=node server-entry.js
directory=/opt/hermes-workspace
environment=
    PORT="3000",
    NODE_ENV="%(ENV_NODE_ENV)s",
    HOST="0.0.0.0",
    HERMES_API_URL="http://localhost:8642",
    HERMES_DASHBOARD_URL="http://localhost:9119",
    HERMES_API_TOKEN="%(ENV_API_SERVER_KEY)s",
    HERMES_DASHBOARD_TOKEN="%(ENV_API_SERVER_KEY)s",
    HERMES_PASSWORD="%(ENV_HERMES_PASSWORD)s",
    HERMES_SECRET_KEY="%(ENV_HERMES_SECRET_KEY)s",
    HERMES_SESSION_SECRET="%(ENV_HERMES_SESSION_SECRET)s",
    HERMES_ALLOW_INSECURE_REMOTE="%(ENV_HERMES_ALLOW_INSECURE_REMOTE)s",
    COOKIE_SECURE="%(ENV_COOKIE_SECURE)s"
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/var/log/supervisor/hermes.log
stderr_logfile=/var/log/supervisor/hermes-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ---------------------------------------------------------------------------
; Aion UI — port 3005
; ---------------------------------------------------------------------------
[program:aionui]
command=bun run webui:prod --remote --no-build
directory=/opt/aionui
environment=AIONUI_PORT="3005",PORT="3005",NODE_ENV="%(ENV_NODE_ENV)s",AIONUI_HOST="0.0.0.0",HOST="0.0.0.0",AIONUI_ALLOW_REMOTE="true",ALLOW_REMOTE="true"
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/var/log/supervisor/aionui.log
stderr_logfile=/var/log/supervisor/aionui-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
