# docker/supervisord.conf
; =============================================================================
; Supervisord — AI Dev Workspace Process Manager
; =============================================================================
; Manages 3 apps in a single container:
;   - OmniRoute  (port 4000)
;   - Hermes     (port 3000)
;   - Aion UI    (port 3005)
; =============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

; ---------------------------------------------------------------------------
; OmniRoute — AI Gateway
; ---------------------------------------------------------------------------
[program:omniroute]
command=npm start
directory=/opt/omniroute
environment=PORT="4000",NODE_ENV="%(ENV_NODE_ENV)s",HOST="0.0.0.0"
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/var/log/supervisor/omniroute.log
stderr_logfile=/var/log/supervisor/omniroute-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ---------------------------------------------------------------------------
; Hermes Workspace
; ---------------------------------------------------------------------------
[program:hermes]
command=node server-entry.js
directory=/opt/hermes
environment=PORT="3000",NODE_ENV="%(ENV_NODE_ENV)s",HOST="0.0.0.0"
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/var/log/supervisor/hermes.log
stderr_logfile=/var/log/supervisor/hermes-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ---------------------------------------------------------------------------
; Aion UI — port 3005
; ---------------------------------------------------------------------------
[program:aionui]
command=bun run webui:prod --remote --no-build
directory=/opt/aionui
environment=AIONUI_PORT="3005",PORT="3005",NODE_ENV="%(ENV_NODE_ENV)s",AIONUI_HOST="0.0.0.0",HOST="0.0.0.0",AIONUI_ALLOW_REMOTE="true",ALLOW_REMOTE="true"
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/var/log/supervisor/aionui.log
stderr_logfile=/var/log/supervisor/aionui-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
