# docker/supervisord.conf
; =============================================================================
; Supervisord — AI Dev Workspace Process Manager
; =============================================================================
; Manages all UI apps and AI services:
;   - OmniRoute         (port 4000)  — LLM router
;   - hermes-gateway    (port 8642)  — OpenAI-compatible API
;   - hermes-dashboard  (port 9119)  — Sessions, memory, skills APIs
;   - Hermes Workspace  (port 3000)  — connects to hermes-agent on localhost:8642
;   - Aion UI           (port 3005)
; =============================================================================

[supervisord]
nodaemon=true
user=hermes
logfile=/opt/data/logs/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/opt/data/logs

; ---------------------------------------------------------------------------
; OmniRoute AI Gateway (port 4000)
; ---------------------------------------------------------------------------
[program:omniroute]
command=npm start
directory=/opt/omniroute
environment=PORT="4000",HOST="0.0.0.0",NODE_ENV="%(ENV_NODE_ENV)s"
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/opt/data/logs/omniroute.log
stderr_logfile=/opt/data/logs/omniroute-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ---------------------------------------------------------------------------
; hermes gateway (port 8642)
; Waits for OmniRoute to be healthy before starting the python API server
; ---------------------------------------------------------------------------
[program:hermes-gateway]
command=bash -c "until curl -s http://localhost:4000 > /dev/null 2>&1; do sleep 2; done; exec hermes gateway run --no-supervise"
directory=/workspace
environment=PATH="/opt/hermes/.venv/bin:/usr/local/bin:%(ENV_PATH)s",HOME="%(ENV_HERMES_HOME)s"
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/opt/data/logs/hermes-gateway.log
stderr_logfile=/opt/data/logs/hermes-gateway.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ---------------------------------------------------------------------------
; hermes dashboard (port 9119)
; ---------------------------------------------------------------------------
[program:hermes-dashboard]
command=bash -c "until curl -s http://localhost:4000 > /dev/null 2>&1; do sleep 2; done; exec hermes dashboard"
directory=/workspace
environment=PATH="/opt/hermes/.venv/bin:/usr/local/bin:%(ENV_PATH)s",HOME="%(ENV_HERMES_HOME)s"
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/opt/data/logs/hermes-dashboard.log
stderr_logfile=/opt/data/logs/hermes-dashboard.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ---------------------------------------------------------------------------
; Hermes Workspace (port 3000)
; ---------------------------------------------------------------------------
[program:hermes-workspace]
command=node server-entry.js
directory=/opt/hermes-workspace
environment=
    PORT="3000",
    NODE_ENV="%(ENV_NODE_ENV)s",
    HOST="0.0.0.0",
    HERMES_API_URL="http://localhost:8642",
    HERMES_DASHBOARD_URL="http://localhost:9119",
    HERMES_API_TOKEN="%(ENV_API_SERVER_KEY)s",
    HERMES_DASHBOARD_TOKEN="%(ENV_API_SERVER_KEY)s",
    HERMES_PASSWORD="%(ENV_HERMES_PASSWORD)s",
    HERMES_SECRET_KEY="%(ENV_HERMES_SECRET_KEY)s",
    HERMES_SESSION_SECRET="%(ENV_HERMES_SESSION_SECRET)s",
    HERMES_ALLOW_INSECURE_REMOTE="%(ENV_HERMES_ALLOW_INSECURE_REMOTE)s",
    COOKIE_SECURE="%(ENV_COOKIE_SECURE)s",
    HERMES_HOME="%(ENV_HERMES_HOME)s"
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/opt/data/logs/hermes-workspace.log
stderr_logfile=/opt/data/logs/hermes-workspace-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ---------------------------------------------------------------------------
; Aion UI (port 3005)
; ---------------------------------------------------------------------------
[program:aionui]
command=bun run webui:prod --remote --no-build
directory=/opt/aionui
environment=AIONUI_PORT="3005",PORT="3005",NODE_ENV="%(ENV_NODE_ENV)s",AIONUI_HOST="0.0.0.0",HOST="0.0.0.0",AIONUI_ALLOW_REMOTE="true",ALLOW_REMOTE="true",OPENAI_API_KEY="%(ENV_OMNIROUTE_API_KEY)s"
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/opt/data/logs/aionui.log
stderr_logfile=/opt/data/logs/aionui-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
